# supervisord.conf
[supervisord]
nodaemon=true

[program:laravel]
command=/usr/local/bin/php /var/www/fullstack/artisan serve --host=0.0.0.0 --port=9000
autostart=true
autorestart=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
